# config/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true
loglevel=debug
user=root

[program:php-fpm]
pidfile=/var/run/php-fpm.pid
command=php-fpm
autostart=true
user=root

[program:nginx]
pidfile=/var/run/nginx.pid
command=/usr/sbin/nginx
autostart=true
user=root
